# etc/cron.d/process_throttler
# cron job for process_throttler
MAILTO=root
*/5 * * * * root /usr/sbin/process_throttler
